# shared/ansible/roles/ca/templates/caserver.service.j2
[Unit]
Description=CA Server

[Service]
User={{ username }}
WorkingDirectory=/home/{{ username }}/ca
ExecStart=gunicorn wsgi:app

[Install]
WantedBy=multi-user.target